# utils/config/vk_postman_start_bot.service
[Unit]
Description=vkPostman Bot

[Service]
ExecStart=/root/projects/vkPostman/.env/bin/python /root/projects/vkPostman/manage.py start_bot
